# config/cron/crontab
SHELL=/bin/bash
 * * * * * borgmatic --verbosity -1 --syslog-verbosity 1 >> /var/log/cron1.log 2>&1

